# distributed/deploy/ansible/roles/zisk-worker/templates/zisk-worker.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=ZisK Worker — distributed proof generation worker
Documentation=https://github.com/0xPolygonHermez/zisk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ zisk_worker_user }}
Group={{ zisk_worker_group }}
WorkingDirectory={{ zisk_worker_work_dir }}

# HOME is needed so the worker resolves ~/.zisk/cache correctly.
Environment=HOME={{ zisk_worker_work_dir }}

{#- Worker arg list — required + conditional appends. Used by both branches. -#}
{%- set worker_args -%}
--config {{ zisk_worker_config_dst }} --proving-key {{ zisk_worker_proving_key }}
{%- if zisk_worker_coordinator_url %} --coordinator-url {{ zisk_worker_coordinator_url }}{% endif -%}
{%- if zisk_worker_id %} --worker-id {{ zisk_worker_id }}{% endif -%}
{%- if zisk_worker_compute_capacity %} --compute-capacity {{ zisk_worker_compute_capacity }}{% endif -%}
{%- if zisk_worker_emulator %} --emulator{% endif -%}
{%- if zisk_worker_asm %} --asm {{ zisk_worker_asm }}{% endif -%}
{%- if zisk_worker_gpu %} --gpu{% endif -%}
{%- if zisk_worker_log_level %} --log-level {{ zisk_worker_log_level }}{% endif -%}
{%- endset -%}
{% if zisk_worker_mpi_enabled -%}
ExecStart={{ zisk_worker_mpirun_bin }} --report-bindings --allow-run-as-root \
    -np {{ zisk_worker_mpi_processes }} \
    -map-by ppr:{{ zisk_worker_mpi_numa_ppr }}:numa \
    --bind-to numa \
    --rank-by slot \
    -x RAYON_NUM_THREADS={{ zisk_worker_mpi_threads }} \
    {{ zisk_worker_binary_dst }} {{ worker_args }}
{%- else %}
ExecStart={{ zisk_worker_binary_dst }} {{ worker_args }}
{%- endif %}
Restart=on-failure
RestartSec=5
LimitNOFILE=65535
Nice=-10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zisk-worker

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ zisk_worker_work_dir }} {{ zisk_worker_log_dir }}
ReadOnlyPaths={{ zisk_worker_config_dir }} {{ zisk_worker_proving_key }}{% if zisk_worker_asm %} {{ zisk_worker_asm }}{% endif %}

[Install]
WantedBy=multi-user.target
