# distributed/deploy/ansible/roles/zisk-coordinator/templates/zisk-coordinator.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=ZisK Coordinator — coordinator server for the ZisK proving system
Documentation=https://github.com/0xPolygonHermez/zisk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ zisk_coordinator_user }}
Group={{ zisk_coordinator_group }}
WorkingDirectory={{ zisk_coordinator_work_dir }}
{#- Coordinator arg list — required + conditional appends. -#}
{%- set coord_args -%}
--config {{ zisk_coordinator_config_dst }} --api-port {{ zisk_coordinator_api_port }}
{%- if zisk_coordinator_cluster_port %} --cluster-port {{ zisk_coordinator_cluster_port }}{% endif -%}
{%- if zisk_coordinator_metrics_port %} --metrics-port {{ zisk_coordinator_metrics_port }}{% endif -%}
{%- if zisk_coordinator_log_level %} --log-level {{ zisk_coordinator_log_level }}{% endif -%}
{%- endset %}
ExecStart={{ zisk_coordinator_binary_dst }} {{ coord_args }}
Restart=on-failure
RestartSec=3
LimitNOFILE=65535
Nice=-10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zisk-coordinator

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ zisk_coordinator_work_dir }}
ReadOnlyPaths={{ zisk_coordinator_config_dir }}

[Install]
WantedBy=multi-user.target
